# tdt/cvs/cdk/root/etc/30-modemswitcher.rules
ACTION=="add", SUBSYSTEM=="usb" , RUN+="/usr/bin/modemctrl.sh %k"

