# tdt/cvs/cdk/root/etc/55-modem.rules
ACTION!="add|change", GOTO="modem_rules_remove"
KERNEL=="ttyUSB[0-9]" RUN+="/usr/bin/modem.sh add %k %s{idVendor} %s{idProduct}"
KERNEL=="ttyACM[0-9]" RUN+="/usr/bin/modem.sh add %k %s{idVendor} %s{idProduct}"
GOTO="modem_rules_end"

LABEL="modem_rules_remove"
ACTION!="remove", GOTO="modem_rules_end"
KERNEL=="ttyUSB[0-9]" RUN+="/usr/bin/modem.sh remove %k %s{idVendor} %s{idProduct}"
KERNEL=="ttyACM[0-9]" RUN+="/usr/bin/modem.sh remove %k %s{idVendor} %s{idProduct}"

LABEL="modem_rules_end"
